# healthd-graph.service
[Unit]
Description=RRD Health Monitor - Graph generation

[Service]
ExecStart=@LIBEXEC@/healthd/healthd-graph

