# healthd-graph.service
[Unit]
Description=RRD Health Monitor - Graph generation

[Service]
ExecStart=/usr/bin/healthd-graph

